# foot/foot-color-monitor.service
[Unit]
Description=Sync foot terminal colors with desktop color scheme (XDG portal)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.dotfiles/foot/foot-color-monitor.sh
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=graphical-session.target
